# systemd/dunst.service
[Service]
ExecStart=/usr/bin/dunst
Restart=always

[Unit]
Description=Dunst
