# systemd/picom.service
[Service]
ExecStart=/usr/bin/picom --experimental-backends
Restart=always

[Unit]
Description=Picom
