# systemd/picom.service
[Service]
ExecStart=/usr/bin/picom
Restart=always

[Unit]
Description=Picom
