# systemd/redshift.service
[Service]
ExecStart=/usr/bin/redshift
Restart=always

[Unit]
Description=Redshift
