# systemd/theme.service
[Service]
ExecStart=%h/.local/bin/theme watch
Restart=always

[Unit]
Description=Theme hot reload
