# systemd/wallpaper.service
[Service]
Type=simple
ExecStart=%h/.local/bin/wallpaper_set

[Unit]
Description=Set theme
